# src/timer.timer
[Unit]
Description=off-time timer

[Timer]
OnCalendar={s}

[Install]
WantedBy=timers.target
